# launchd/com.rout.proactive-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rout.proactive-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>scripts/proactive_agent.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>INSTALL_DIR</string>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/REPLACE_WITH_YOUR_USERNAME/.openclaw/logs/proactive-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_WITH_YOUR_USERNAME/.openclaw/logs/proactive-agent.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/REPLACE_WITH_YOUR_USERNAME</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
