# launchd/com.rout.status-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rout.status-api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/__USERNAME__/.openclaw/workspace/rout/monitor/status_api.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/__USERNAME__/.openclaw/logs/status-api.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/__USERNAME__/.openclaw/logs/status-api.err</string>

    <key>WorkingDirectory</key>
    <string>/Users/__USERNAME__/.openclaw/workspace/rout</string>
</dict>
</plist>
